# .service
[Unit]
Description=Mesh Network Master
After=network.target

StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
WorkingDirectory=/srv/git/mesh
Environment="PATH=/srv/git/mesh/env/bin"
ExecStart=/srv/git/server/env/bin/python3 -m main.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
